# celery_flower.ini
[program:celery-flower]
command=/home/git/iclassvideo/.env/bin/celery flower -A iclassvideo.celery_app
directory=/home/git/
autostart=true
autorestart=true
priority=120
user=root
redirect_stderr=true
stdout_logfile=/var/log/supervisor/celery-flower.log
stdout_logfile_maxbytes=500MB
