# celery_record2.ini
[program:celery-record2]
command=/home/git/iclassvideo/.env/bin/celery -A iclassvideo.celery_app worker -Q ffmpeg_record_tasks  -c100 -Ofair  -n recordworker2@%%n -l info
directory=/home/git
autostart=true
autorestart=true
priority=99
user=root
redirect_stderr=true
stdout_logfile=/var/log/supervisor/celery-record.log
stdout_logfile_maxbytes=500MB
