# celery_fmt.ini
[program:celery-fmt]
command=/usr/local/bin/celery -A celery_task.celery_app worker -Q ffmpeg_format_tasks -c4 -Ofair -n formatworker@%%n -l info
directory=/home/git
autostart=true
autorestart=true
priority=99
user=root
redirect_stderr=true
stdout_logfile=/var/log/supervisor/celery-fmt.log
stdout_logfile_maxbytes=500MB
