# celery_controller.ini
[program:celery-controller]
command=/home/git/iclassvideo/.env/bin/celery -A iclassvideo.celery_app worker -Q issue_record_tasks,issue_format_tasks -n controller@%%n -B -s /opt/iclassvideo/celery_controller_beat -l info
directory=/home/git/
autostart=true
autorestart=true
priority=99
user=root
redirect_stderr=true
stdout_logfile=/var/log/supervisor/celery-controller.log
stdout_logfile_maxbytes=500MB
