# fsmon.ini
[program:fsmon]
command=/home/git/iclassvideo/.env/bin/python iclassvideo/fsmon.py
directory=/home/git/
autostart=true
autorestart=true
priority=99
user=root
redirect_stderr=true
stdout_logfile=/var/log/supervisor/fsmon.log
stdout_logfile_maxbytes=500MB
